# ollama_app.service
[Unit]
Description=Simple Python Server
After=network.target

[Service]
WorkingDirectory=/home/ubuntu/Ollama
ExecStart=/home/ubuntu/Ollama/.venv/bin/python /home/ubuntu/Ollama/app.py
Restart=always

[Install]
WantedBy=multi-user.target