# ollama_app.service
[Unit]
Description=Ollama Quiz App Service
After=network.target

[Service]
User=ameertabri
WorkingDirectory=/home/ameertabri/OllamaProject
Environment="PATH=/home/ameertabri/OllamaProject/venv/bin"
ExecStart=/home/ameertabri/OllamaProject/venv/bin/python backend/main.py
Restart=always

[Install]
WantedBy=multi-user.target 